# systemd/modemproxy-bandwidth.timer
[Unit]
Description=Sample modemproxy bandwidth every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5

[Install]
WantedBy=timers.target
